# profile.d/zfs_grub.sh
export ZPOOL_VDEV_NAME_PATH=1
ZPOOL_VDEV_NAME_PATH=1
